# etc/systemd/system/btrfs-snapshot-and-backup.service
[Unit]
Description=Make BTRFS snapshots and copy them

[Service]
ExecStart=/home/lucas/Documentos/scripts/btrbk_
Type=oneshot
