# etc/systemd/system/check_btrfs_health.timer
[Unit]
Description=Check BTRFS storage devices' health

[Timer]
OnCalendar=*:0,30

[Install]
WantedBy=timers.target
