# etc/systemd/system/sync-repos.service
[Unit]
Description=Sync repos

[Service]
ExecStart=/usr/sbin/emaint sync -A
Type=oneshot
