# etc/systemd/system/winetricks_update.timer
[Unit]
Description=Run winetricks update script weekly (Saturday)

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
