# etc/systemd/system/winetricks_update.service
[Unit]
Description=Run winetricks update script
After=network.target

[Service]
ExecStart=/usr/bin/winetricks --self-update
Type=oneshot
